# src/etc/services.d/xbs-api/run
#!/usr/bin/with-contenv sh

cd /xbs-api

exec s6-setuidgid xbs node dist/api.js
